# tempmon_supervisord.conf
[program:tempmon_livingroom]
command=python /home/jdsnape/workers/tempmon/tempmon.py -c /home/jdsnape/workers/tempmon/livingroom.conf
